# docker/root/etc/services.d/rix_template/run
#!/usr/bin/with-contenv bash

exec rix_template --port=8080 --config=/config --example=$EXAMPLE
